# supervisor/earnfm.conf
[program:earnfm]
command=/opt/earnfm/earnfm
user=root
stdout_logfile=/home/logs/earnfm.log
stderr_logfile=/home/logs/earnfm.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
priority=4