# supervisor/earnapp.conf
[program:earnapp]
command=install
user=root
stdout_logfile=/home/logs/earnapp.log
stderr_logfile=/home/logs/earnapp.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
priority=3
