# supervisor/meson.conf
[program:meson]
command=/opt/meson/meson.sh
user=root
stdout_logfile=/home/logs/meson.log
stderr_logfile=/home/logs/meson.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
priority=8
